# files/etc/crontabs/root
00 3 * * 0 sleep 5 && touch /etc/banner && reboot
